# sh/profile.d/cd.sh
# change to, and immediately list, a directory
cdls() { cd "$@" && ls; }

# change to, and immediately list (at length), a directory
cdll() { cd "$@" && ll; }
