# sh/profile.d/cp.sh
cp() { command cp -aiv "$@"; }
#                  ││└─ verbose
#                  │└── interactive
#                  └─── archive mode (recursive; don't follow symlinks; 
#                       preserve attributes)
