# sh/profile.d/dtrx.sh
if _inPath dtrx; then
  dtrx() { command dtrx --verbose "$@"; }
fi
