# sh/profile.d/inkscape.sh
if command -v inkscape >/dev/null; then
  export INKSCAPE_PROFILE_DIR="$HOME/share/inkscape"
fi
