# sh/profile.d/ln.sh
ln() { command ln -v "$@"; }
#                  └─ verbose
