# sh/profile.d/pip.sh
if command -v pip >/dev/null; then
  export PIP_CONFIG_FILE="$HOME/.local/config/pip.conf"
fi
