# sh/profile.d/rename.sh
if _inPath rename; then
  rename() { command rename --verbose "$@"; }
fi
