# sh/profile.d/rm.sh
rm() { command rm -iv "$@"; }
#                  │└─ verbose
#                  └── interactive
