# sh/profile.d/scp.sh
scp() { command scp -rp "$@"; }
#                    │└── preserve times/modes
#                    └─── recursive
