# sh/profile.d/stow.sh
if _inPath stow; then
  stow() { command stow --verbose "$@"; }
fi
