# sh/profile.d/unar.sh
if _inPath unar; then
  unar() { command unar -copy-time -forks hidden "$@"; }
fi
