# sh/profile.d/xterm.sh
if [ x"$XTERM_VERSION" != "x" ]; then
  export XTERM_PROFILE="yes"
fi
